# docker-manifest/mautic.crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 4 15 * *  root  php /var/www/html/app/console mautic:iplookup:download  > /proc/1/fd/1 2> /proc/1/fd/2
* * *  * *  root  echo "[$(date --iso-8601=ns)] Tick...  Tock..."         > /proc/1/fd/1 2> /proc/1/fd/2